# build/logger/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:log-generator]
command=/usr/local/bin/log-generator /etc/log-generator/conf/apache_2_4_access_log.conf

[program:filebeat]
command=/usr/share/filebeat/bin/filebeat -e
    -c /etc/filebeat/filebeat.yml 
    -path.home /usr/share/filebeat
    -path.config /etc/filebeat
    -path.data /var/lib/filebeat
    -path.logs /var/log/filebeat
